# docker/build/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/usr/sbin/sshd -D

[program:apache2]
command=apachectl -D FOREGROUND