# docker/build/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/root/run.sh

[program:apache2]
command=apachectl -D FOREGROUND